# tranquility/supervisord.conf
[supervisord]
nodaemon=true
loglevel=debug


[program:druid-tranquility]
user=druid
command=/var/lib/druid/starttranquility.sh %(ENV_ZOOKEEPER_HOST)s
redirect_stderr=true
priority=100


